# deb/DEBIAN/postinst
#!/bin/bash

mkdir -m 0666 /usr/share/wall
systemctl daemon-reload
systemctl enable wallpaperloader.service
systemctl start wallpaperloader

