# deb/DEBIAN/prerm
#!/bin/bash

systemctl stop wallpaperloader
systemctl disable wallpaperloader
exit 0

